# rootfs/etc/services.d/mariadb/run
#!/usr/bin/with-contenv sh

exec /usr/bin/mysqld_safe --user=daspanel --console --port=3306


